# snapdog/systemd/snapdog.service
[Unit]
Description=SnapDog multiroom audio server
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/snapdog
ExecStart=/usr/bin/snapdog $SNAPDOG_ARGS
Restart=on-failure
RestartSec=5
DynamicUser=yes
SupplementaryGroups=audio
StateDirectory=snapdog
WorkingDirectory=/var/lib/snapdog
ConfigurationDirectory=snapdog
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
